# etc/profile.d/gistow.sh
# Should be a better choice instead of the script
alias gistow='git --work-tree .. --git-dir'
